# sixad/98-sixad.rules
# do not edit this file, it will be overwritten on update

ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/sixad --udev"
# Please note that there are name filters on the 'sixad' program
# sixad is only started if the last detected device has name 'PLAYSTATION(R)3 Controller' or 'Wireless Keypad'
